# src/etc/crontabs/root
0 0 * * 1 rm /etc/smartrouter/reconnect.log
*/5 * * * * sh /etc/smartrouter/checkpoint.sh
0 5 * * * sh /etc/smartrouter/x-reboot.sh
